# home/.config/systemd/user/blinkymeter.service
[Unit]
Description=Blinky Meter

[Service]
ExecStart=%h/bin/blinkymeter

[Install]
WantedBy=default.target
